# config/cluster2/ubuntu.arm64v8/master-pi/etc/systemd/system/snap-lxd-24326.mount
[Unit]
Description=Mount unit for lxd, revision 24326
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/lxd_24326.snap
Where=/snap/lxd/24326
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
