# ansible/cluster1/mgt/roles/mount_disks/templates/mnt-disks-vol3.mount
[Unit]
Description = Mount for Kubernetes vol3

[Mount]
What=/dev/mmcblk0p7
Where=/mnt/disks/vol3
Type=ext4

[Install]
WantedBy = multi-user.target
